# packages/extra/orionos-containers/systemd/orionos-containers.service
[Unit]
Description=OrionOS Container Runtime Setup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/orionos-containers-setup setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
