# packages/extra/orionos-cloud-sync/systemd/orionos-cloud-sync.service
[Unit]
Description=OrionOS Cloud Sync Service
After=graphical-session.target network-online.target
Wants=graphical-session.target network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/orionos-cloud-sync-launcher --daemon
Restart=on-failure
RestartSec=30
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
