# packages/extra/orionos-display-optimization/systemd/orionos-display-optimized.service
[Unit]
Description=OrionOS Display Optimization Daemon
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/python3/site-packages/orionos_display_optimization.py --daemon
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
