# packages/core/orionos-voice-assistant/systemd/orionos-voice-assistant.service
[Unit]
Description=OrionOS Voice Assistant
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/orionos-voice-assistant-launcher
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
